# root/etc/services.d/nfs/run
#!/usr/bin/with-contenv bash

mkdir -p "$MOUNTPOINT"
/sbin/rpcbind -f &
/bin/mount -t "$FSTYPE" -o "$MOUNT_OPTIONS" "$SERVER:$SHARE" "$MOUNTPOINT"